# enterprise-monitor.service
[Service]
Type=simple
User=enterprise-monitor
Group=enterprise-monitor
WorkingDirectory=/opt/enterprise-monitor
Environment=PATH=/opt/enterprise-monitor/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/enterprise-monitor/venv/bin/python server_launcher.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=DISPLAY=:99
Environment=ENTERPRISE_MONITOR_SERVER=true
Environment=CHROME_NO_SANDBOX=true

# Giảm bớt security restrictions
NoNewPrivileges=false
# PrivateTmp=true
# ProtectHome=true
# ProtectSystem=strict
ReadWritePaths=/opt/enterprise-monitor /var/log/enterprise-monitor /tmp /home

# Resource limits
MemoryMax=2G
CPUQuota=80%
[Install]
WantedBy=multi-user.target